# dot_config/systemd/user/blurt.service.d/override.conf
# Recover from *hangs*, not just crashes. The 2026-07-24 incident: daemon ran
# 3d+, wedged (0.1s recordings, transcripts 14-49s late) but never exited, so
# Restart=on-failure couldn't catch it. Restart=always covers clean-exit stalls;
# RuntimeMaxSec=1d forces a daily recycle (systemd stops after 24h uptime,
# Restart=always brings it back in ~4s) so state never rots.
[Service]
Restart=always
RuntimeMaxSec=1d
